# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

; GPU worker owns the model children; give it time to finish a page on shutdown.
[program:worker]
command=/opt/venv/bin/python -m worker.run
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=5
stopsignal=TERM
stopwaitsecs=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:server]
command=node /app/apps/server/dist/index.mjs
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=3
stopsignal=TERM
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
